# backend/parsers/drone/run_as_service/drone.service
[Unit]
Description=Drone Parser ROS Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__WORKING_DIR__
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=__WORKING_DIR__/.env
ExecStart=/bin/bash -lc '__WORKING_DIR__/start_drone.sh'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

